# deploy/secretary-ai.service
[Unit]
Description=Secretary AI Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=secretary
WorkingDirectory=/home/secretary/secretary-ai
EnvironmentFile=/home/secretary/secretary-ai/.env
ExecStart=/home/secretary/secretary-ai/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
