# deploy/secretary-ai.service
[Unit]
Description=Secretary AI Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=murad
WorkingDirectory=/home/murad/secretary-ai
EnvironmentFile=/home/murad/secretary-ai/.env
ExecStart=/home/murad/secretary-ai/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
